# trapdoor/trapdoor.service
[Unit]
Description=Trapdoor Countdown Service
After=network.target

[Service]
User=REPLACE_WITH_USER
WorkingDirectory=/opt/trapdoor
ExecStart=/opt/trapdoor/venv/bin/gunicorn -w 1 -b 127.0.0.1:5001 app:app
Restart=always
RestartSec=5
Environment="PATH=/opt/trapdoor/venv/bin"

[Install]
WantedBy=multi-user.target
